# api_site.service
[Unit]
Description=api_site

[Service]
Type=simple
ExecStart=/etc/init.d/api_site start

[Install]
WantedBy=multi-user.target
